# AVS/debug/android-29/armeabi-v7a/nghttp2/contrib/nghttpx-upstart.conf
# vim: ft=upstart:

description "HTTP/2 reverse proxy"

start on runlevel [2]
stop on runlevel [016]

exec /home2/silogood/alexa/AVS/debug/android-29/armeabi-v7a/install/bin/nghttpx
